# deploy/remotex-daemon.service
[Unit]
Description=Remotex daemon (bridges codex app-server to the relay)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@@WORKING_DIR@@
Environment=PATH=@@VENV_BIN@@:/usr/local/bin:/usr/bin:/bin
ExecStart=@@VENV_BIN@@/python3 -m daemon run --config @@CONFIG_PATH@@
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=default.target
